# runit/etc/sv/denyhosts/finish
#!/bin/sh
ps aux | grep denyhosts | grep foreground | grep -v grep | awk '{print $2}' | xargs kill
