# runit/etc/sv/denyhosts/run
#!/bin/sh
exec /usr/sbin/denyhosts --noemail --foreground --config=/etc/denyhosts.conf
